# koha/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root
; umask=000 lets starman create the plack Unix socket as 0777 (world-writable)
; so Apache (www-data) can proxy requests to it without permission errors
umask=000

[program:apache2]
; Apache runs in foreground — supervisord owns it directly
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stderr_logfile=/var/log/apache2/error.log
stdout_logfile=/var/log/apache2/access.log
priority=10

[program:zebra]
; koha-zebra --start launches a self-managed daemon (daemon --respawn wrapper)
; autorestart=false: supervisord starts it once; daemon wrapper handles respawns
command=koha-zebra --start %(ENV_KOHA_INSTANCE)s
autostart=true
autorestart=false
startsecs=0
exitcodes=0,1
stderr_logfile=/var/log/koha/%(ENV_KOHA_INSTANCE)s/zebra.log
stdout_logfile=/var/log/koha/%(ENV_KOHA_INSTANCE)s/zebra.log
priority=15

[program:koha-plack]
; koha-plack-wrapper sets umask 0 before starting so plack.sock is 0777 (Apache-accessible)
command=koha-plack-wrapper %(ENV_KOHA_INSTANCE)s
autostart=true
autorestart=false
startsecs=0
exitcodes=0,1
stderr_logfile=/var/log/koha/%(ENV_KOHA_INSTANCE)s/plack.log
stdout_logfile=/var/log/koha/%(ENV_KOHA_INSTANCE)s/plack.log
priority=20

[program:koha-sip]
; koha-sip --start launches SIP server as a daemon
command=koha-sip --start %(ENV_KOHA_INSTANCE)s
autostart=true
autorestart=false
startsecs=0
exitcodes=0,1
stderr_logfile=/var/log/koha/%(ENV_KOHA_INSTANCE)s/sip.log
stdout_logfile=/var/log/koha/%(ENV_KOHA_INSTANCE)s/sip.log
priority=25
